# scripts/thunder-ptm-service.service
[Unit]
Description=Thunder PTM Service
After=network.target

[Service]
Restart=always
User=thunder-ptm
Group=thunder-ptm
Environment=NODE_ENV=production
Environment=NODE_VERSION=node
WorkingDirectory=/home/thunder-ptm/thunder-performance-task-manager
ExecStart=/home/thunder-ptm/.nvm/nvm-exec npm run service

[Install]
WantedBy=multi-user.target
